# openAI_Agentic_SOC_Analyst/ctf_scenarios/rdp_password_spray.json
{
  "scenario_name": "RDP Password Spray - Full Compromise",
  "scenario_id": "rdp_password_spray",
  "description": "Suspicious RDP login activity detected on cloud-hosted Windows server. Multiple failed attempts followed by successful login.",
  "device_filter": "where DeviceName contains 'flare'",
  "incident_date": "2025-09-14",
  "total_flags": 10,
  
  "flags": {
    "1": {
      "title": "Attacker IP Address",
      "stage": "Initial Access",
      "mitre": "T1110.001 - Brute Force: Password Guessing",
      "objective": "Identify the external IP address that successfully logged in via RDP after a series of failures",
      "scenario_context": "Suspicious RDP login activity has been detected. Multiple failed attempts were followed by a successful login, suggesting brute-force or password spraying behavior.",
      "guidance": "Review the authentication telemetry and look for signs of repeated failed logins followed by a successful one. Focus on logins that originated from external IP addresses.",
      "table": "DeviceLogonEvents",
      "key_fields": ["RemoteIP", "ActionType", "RemoteIPType", "AccountName", "Timestamp"],
      "flag_format": "xxx.xxx.xxx.xxx",
      "answer": "159.26.106.84",
      "hints": [
        "Filter by RemoteIPType == 'Public'",
        "Look for ActionType == 'LogonSuccess'",
        "Sort by Timestamp to find earliest",
        "Device filter: DeviceName contains 'flare'"
      ]
    },
    
    "2": {
      "title": "Compromised Account",
      "stage": "Initial Access",
      "mitre": "T1078 - Valid Accounts",
      "objective": "Determine the username that was used during the successful RDP login associated with the attacker's IP",
      "guidance": "Pivot from the successful login identified in Flag 1. Analyze the associated account used in that authentication event.",
      "table": "DeviceLogonEvents",
      "key_fields": ["AccountName", "ActionType", "RemoteIP"],
      "flag_format": "username",
      "answer": "slflare",
      "depends_on": [1],
      "correlation": {
        "uses_flag_1": "attacker_ip",
        "filter_hint": "where RemoteIP == '{flag_1_answer}'"
      },
      "hints": [
        "Use the IP from Flag 1 as filter",
        "Look for LogonSuccess with that IP",
        "Extract AccountName field"
      ]
    },
    
    "3": {
      "title": "Executed Binary Name",
      "stage": "Execution",
      "mitre": "T1059.003 - Command and Scripting Interpreter",
      "objective": "Identify the name of the binary executed by the attacker after gaining RDP access",
      "guidance": "Focus your investigation on process execution under the compromised user account from Flag 2. Look for binaries launched from unusual paths like Public, Temp, or Downloads folders.",
      "table": "DeviceProcessEvents",
      "key_fields": ["FileName", "ProcessCommandLine", "AccountName", "Timestamp"],
      "flag_format": "filename",
      "answer": "msupdate.exe",
      "depends_on": [2],
      "correlation": {
        "uses_flag_2": "compromised_account",
        "filter_hint": "where AccountName == '{flag_2_answer}'"
      },
      "hints": [
        "Filter by AccountName from Flag 2",
        "Look for ProcessCommandLine containing 'Public', 'Temp', or 'Bypass'",
        "Find earliest suspicious binary execution"
      ]
    },
    
    "4": {
      "title": "Command Line Used to Execute Binary",
      "stage": "Execution",
      "mitre": "T1059 - Command and Scripting Interpreter",
      "objective": "Provide the full command line used to launch the binary from Flag 3",
      "guidance": "Review command-line arguments associated with process execution under the compromised account. Pay attention to how the attacker invoked the binary and any parameters used.",
      "table": "DeviceProcessEvents",
      "key_fields": ["ProcessCommandLine", "FileName", "AccountName"],
      "flag_format": "Full command line, including the path",
      "answer": "\"msupdate.exe\" -ExecutionPolicy Bypass -File C:\\Users\\Public\\update_check.ps1",
      "depends_on": [2, 3],
      "correlation": {
        "uses_flag_2": "compromised_account",
        "uses_flag_3": "binary_name",
        "filter_hint": "where FileName == '{flag_3_answer}' and AccountName == '{flag_2_answer}'"
      }
    },
    
    "5": {
      "title": "Persistence Mechanism Created",
      "stage": "Persistence",
      "mitre": "T1053.005 - Scheduled Task/Job: Scheduled Task",
      "objective": "Identify the name of the scheduled task created by the attacker",
      "guidance": "Investigate task creation activity by examining both process execution events and registry modifications. Look for newly registered scheduled tasks in the Windows TaskCache registry.",
      "table": "DeviceRegistryEvents",
      "key_fields": ["RegistryKey", "RegistryValueName", "InitiatingProcessAccountName", "Timestamp"],
      "flag_format": "Name of the scheduled task",
      "answer": "MicrosoftUpdateSync",
      "depends_on": [2],
      "correlation": {
        "uses_flag_2": "compromised_account",
        "filter_hint": "where RegistryKey contains 'TaskCache' and RegistryKey contains 'Tree'"
      }
    },
    
    "6": {
      "title": "Defender Setting Modified",
      "stage": "Defense Evasion",
      "mitre": "T1562.001 - Impair Defenses: Disable or Modify Windows Defender",
      "objective": "Identify the folder path that was excluded from Defender scans",
      "guidance": "Look for registry modification events linked to Defender exclusions. Focus on exclusions that prevent scanning of specific folders.",
      "table": "DeviceRegistryEvents",
      "key_fields": ["RegistryKey", "RegistryValueName", "RegistryValueData", "Timestamp"],
      "flag_format": "Full folder path (e.g., C:\\Example\\Path)",
      "answer": "C:\\Windows\\Temp",
      "depends_on": [2],
      "correlation": {
        "filter_hint": "where RegistryKey contains 'Defender' and RegistryKey contains 'Exclusions'"
      }
    },
    
    "7": {
      "title": "Discovery Command Run",
      "stage": "Discovery",
      "mitre": "T1082 - System Information Discovery",
      "objective": "Identify the exact command line the attacker executed to perform system discovery",
      "guidance": "Review process execution data for evidence of built-in Windows tools used for enumeration. Pay attention to commands like systeminfo, whoami, ipconfig, etc.",
      "table": "DeviceProcessEvents",
      "key_fields": ["ProcessCommandLine", "FileName", "AccountName", "Timestamp"],
      "flag_format": "The full command line, including the path",
      "answer": "\"cmd.exe\" /c systeminfo",
      "depends_on": [2],
      "correlation": {
        "uses_flag_2": "compromised_account",
        "filter_hint": "where AccountName == '{flag_2_answer}' and ProcessCommandLine contains 'systeminfo' or ProcessCommandLine contains 'whoami'"
      }
    },
    
    "8": {
      "title": "Archive File Created",
      "stage": "Collection",
      "mitre": "T1560.001 - Archive Collected Data: Local Archiving",
      "objective": "Identify the name of the archive file created by the attacker",
      "guidance": "Look for file creation or process activity involving archiving tools. Focus on .zip, .rar, or .7z files created in non-standard directories.",
      "table": "DeviceProcessEvents",
      "key_fields": ["ProcessCommandLine", "FileName", "Timestamp"],
      "flag_format": "Archive filename",
      "answer": "backup_sync.zip",
      "depends_on": [2],
      "correlation": {
        "uses_flag_2": "compromised_account",
        "filter_hint": "where ProcessCommandLine contains '.zip' or ProcessCommandLine contains '.rar'"
      }
    },
    
    "9": {
      "title": "C2 Connection Destination",
      "stage": "Command and Control",
      "mitre": "T1071.001 - Application Layer Protocol: Web Protocols",
      "objective": "Identify the destination the attacker's C2 beacon connected to or retrieved tooling from",
      "guidance": "Review outbound network connections tied to suspicious activity. Look for external IPs or domains contacted shortly after initial execution or persistence.",
      "table": "DeviceProcessEvents",
      "key_fields": ["ProcessCommandLine", "RemoteIP", "RemoteUrl"],
      "flag_format": "IP or domain",
      "answer": "185.92.220.87",
      "depends_on": [2],
      "correlation": {
        "filter_hint": "where ProcessCommandLine contains 'http' and ProcessCommandLine contains '.zip'"
      }
    },
    
    "10": {
      "title": "Exfiltration Attempt Detected",
      "stage": "Exfiltration",
      "mitre": "T1048.003 - Exfiltration Over Unencrypted Protocol",
      "objective": "Identify the external IP address and port used during data exfiltration attempt",
      "guidance": "Search DeviceNetworkEvents for outbound traffic occurring after the archive was created. Focus on unusual external destinations.",
      "table": "DeviceProcessEvents",
      "key_fields": ["ProcessCommandLine", "RemoteIP", "RemotePort"],
      "flag_format": "IP:Port",
      "answer": "185.92.220.87:8081",
      "depends_on": [8, 9],
      "correlation": {
        "uses_flag_9": "c2_ip",
        "filter_hint": "Look for the C2 IP from Flag 9 with port number"
      }
    }
  },
  
  "attack_chain": [
    {"stage": "Initial Access", "flags": [1, 2], "description": "External brute-force (159.26.106.84) → Account compromise (slflare)"},
    {"stage": "Execution", "flags": [3, 4], "description": "Malicious binary execution (msupdate.exe) with PowerShell bypass"},
    {"stage": "Persistence", "flags": [5], "description": "Scheduled task creation (MicrosoftUpdateSync)"},
    {"stage": "Defense Evasion", "flags": [6], "description": "Defender exclusion added (C:\\Windows\\Temp)"},
    {"stage": "Discovery", "flags": [7], "description": "System enumeration (systeminfo command)"},
    {"stage": "Collection", "flags": [8], "description": "Data staging (backup_sync.zip archive)"},
    {"stage": "Command & Control", "flags": [9], "description": "C2 communication (185.92.220.87)"},
    {"stage": "Exfiltration", "flags": [10], "description": "Data exfiltration attempt (185.92.220.87:8081)"}
  ]
}

